# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:api]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=1
stdout_events_enabled = true
stderr_events_enabled = true

[program:agent]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=100 
stdout_events_enabled = true
stderr_events_enabled = true

[program:aggregator]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:alarm]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10
stdout_events_enabled = true
stderr_events_enabled = true

[program:gateway]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:hbs]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:judge]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:nodata]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:transfer]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=10 
stdout_events_enabled = true
stderr_events_enabled = true

[program:graph]
user = open-falcon
process_name=%(program_name)s
directory = /usr/local/open-falcon/
command = /usr/local/open-falcon/bin/falcon-%(program_name)s -c /usr/local/open-falcon/config/%(program_name)s-cfg.json
autorestart=unexpected
autostart=true
startretries=5
stopwaitsecs=10
priority=1000 
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler